# deploy/docker/allinone/apiserver.service
[Unit]
Description=MCP API Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/apiserver
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target