# deploy/docker/allinone/mock-user-svc.service
[Unit]
Description=MCP Mock User Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mock-user-svc
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 