# deploy/docker/allinone/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:mcp-gateway]
command=/usr/local/bin/mcp-gateway
directory=/app
autostart=true
autorestart=true
stderr_logfile=NONE
stdout_logfile=NONE

[program:mock-user-svc]
command=/usr/local/bin/mock-user-svc
directory=/app
autostart=true
autorestart=true
stderr_logfile=NONE
stdout_logfile=NONE

[program:apiserver]
command=/usr/local/bin/apiserver
directory=/app
autostart=true
autorestart=true
stderr_logfile=NONE
stdout_logfile=NONE

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=NONE
stderr_logfile=NONE